# archiso/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -mx -- \\u' --noclear --autologin root %I $TERM
